# gateway/patient-care-gateway.service
# Patient Care gateway — Arduino serial <-> backend bridge, run on the Raspberry Pi.
#
# Install:
#   sudo cp gateway/patient-care-gateway.service /etc/systemd/system/
#   sudo systemctl daemon-reload
# Enable + start ONLY AFTER you have pasted a DEVICE_TOKEN into gateway/.env
# (the process exits immediately with code 2 if DEVICE_TOKEN is empty):
#   sudo systemctl enable --now patient-care-gateway
# Logs:
#   journalctl -u patient-care-gateway -f

[Unit]
Description=Patient Care Gateway (Arduino serial <-> backend)
After=network-online.target
Wants=network-online.target
# Never stop trying to (re)start on a cold boot: if the Arduino or the WiFi
# isn't ready yet, the process may exit a few times in quick succession.
# Disable systemd's start-rate limit so it keeps retrying forever.
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
# Grant access to /dev/ttyACM0 (owned by group dialout).
SupplementaryGroups=dialout
WorkingDirectory=/home/pi/Downloads/Patient-Care-main/gateway
# config.py uses python-dotenv to read ./ .env from this WorkingDirectory.
ExecStart=/home/pi/Downloads/Patient-Care-main/gateway/.venv/bin/python gateway.py
# always (not on-failure): recover from any exit — serial not ready, transient
# crash, clean exit. The backend WS itself reconnects forever inside the process.
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
